# build/service/supervisord.conf
[supervisord]
nodaemon=true
user=root

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[unix_http_server]
file = /tmp/supervisor.sock
chmod = 0777
chown= root:root
username=admin
password=default

[supervisorctl]
serverurl = unix:///tmp/supervisor.sock
username=admin
password=default

[program:gunicorn]
command=/usr/local/bin/gunicorn --config /etc/gunicorn/conf.d/gunicorn.conf -b unix:/tmp/gunicorn.sock service.app:app
directory=/service
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:nginx]
command=/usr/sbin/nginx -g 'daemon off;'
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
